# monitor_app/crontab
* * * * * python /Users/xavier/Documents/Code\ Project/monitor_load_project/monitor_app/manage.py update_data
